# services/openclaw-config/auto-brain-ingest.service
[Unit]
Description=Auto Brain Ingest - Watches MEMORY.md and updates brain.db
After=network.target

[Service]
Type=simple
User=jetson
WorkingDirectory=/home/jetson/.openclaw/workspace
ExecStart=/usr/bin/python3 /home/jetson/.openclaw/workspace/auto_brain_ingest.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
