# services/openclaw-config/periodic-brain-ingest.service
[Unit]
Description=Periodic Brain Ingest - runs every minute

[Service]
Type=oneshot
User=jetson
WorkingDirectory=/home/jetson/.openclaw/workspace
ExecStart=/home/jetson/.openclaw/workspace/cron_ingest.sh

[Install]
WantedBy=multi-user.target
